# deploy/roadwatch-ai.service
[Unit]
Description=RoadWatch AI traffic monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=roadwatch
Group=roadwatch
WorkingDirectory=/opt/roadwatch-ai
EnvironmentFile=/opt/roadwatch-ai/.env
ExecStart=/opt/roadwatch-ai/.venv/bin/python -m roadwatch_ai run --config /opt/roadwatch-ai/config.yaml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
